# crontab
* * * * * /bin/bash /opt/machine-system/stats.sh >/dev/null 2>&1
